# etc/xdg/autostart/display.desktop
[Desktop Entry]
Name=OutFox
Exec=/home/pi/OutFox/OutFox

